# deb/DEBIAN/postinst
#!/bin/sh
update-rc.d ringd defaults 2>/dev/null
exit 0
