# deb/DEBIAN/prerm
#!/bin/sh
update-rc.d -f ringd remove > /dev/null
exit 0
